# movie/conf/cron/crontab
*/30 * * * * date
*/30 * * * * /download_current_rss.sh